# config/etc/systemd/system/frr.service.d/after-vpp.conf
[Unit]
# Ensure FRR starts after VPP has created TAP interfaces
After=vpp-core-config.service
Wants=vpp-core-config.service
