# config/templates/systemd/vpp-module.service.j2
# VPP {{ module.display_name }} Module Service
# Generated by configure-router.py

[Unit]
Description=VPP {{ module.display_name }} Module
After=vpp-core.service
Requires=vpp-core.service
ConditionPathExists=/etc/vpp/startup-{{ module.name }}.conf

[Service]
Type=simple
ExecStart=/usr/bin/vpp -c /etc/vpp/startup-{{ module.name }}.conf
ExecStopPost=/bin/rm -f {{ module.socket_path }}
Restart=on-failure
RestartSec=5
NetworkNamespacePath=/run/netns/dataplane

[Install]
WantedBy=multi-user.target
